# src/ClusterBootstrap/template/cloud-config/infra.upgrade.list
./deploy/ssl/apiserver,/etc/kubernetes/ssl
./deploy/ssl/aggregator/ca.crt,/etc/kubernetes/pki/ca.crt
./deploy/ssl/aggregator/ca.key,/etc/kubernetes/pki/ca.key
./deploy/ssl/aggregator/proxy-client.crt,/etc/kubernetes/pki/proxy-client.crt
./deploy/ssl/aggregator/proxy-client.key,/etc/kubernetes/pki/proxy-client.key
./deploy/master/basicauth,/etc/kubernetes/basicauth
./deploy/master/worker-kubeconfig.yaml,/etc/kubernetes/worker-kubeconfig.yaml
./deploy/master/kube-apiserver.yaml,/etc/kubernetes/manifests/kube-apiserver.yaml
./deploy/master/kube-controller-manager.yaml,/etc/kubernetes/manifests/kube-controller-manager.yaml
./deploy/master/kube-scheduler.yaml,/etc/kubernetes/manifests/kube-scheduler.yaml
./deploy/bin/kubelet,/opt/bin/kubelet
./deploy/bin/kubectl,/opt/bin/kubectl
./deploy/kube-addons,/opt/addons/kube-addons
./deploy/master/dns-kubeconfig.yaml,/etc/kubernetes/dns-kubeconfig.yaml
./deploy/services/clusterroles/clusterrolebindings.yaml,/etc/kubernetes/clusterroles/clusterrolebindings.yaml
./deploy/services/clusterroles/clusterroles.yaml,/etc/kubernetes/clusterroles/clusterroles.yaml
./deploy/bin/macvlan,/opt/cni/bin/macvlan
./deploy/bin/dhcp,/opt/cni/bin/dhcp
./deploy/bin/loopback,/opt/cni/bin/loopback
./deploy/bin/ptp,/opt/cni/bin/ptp
./deploy/bin/ipvlan,/opt/cni/bin/ipvlan
./deploy/bin/bridge,/opt/cni/bin/bridge
./deploy/bin/tuning,/opt/cni/bin/tuning
./deploy/bin/host-local,/opt/cni/bin/host-local
./deploy/bin/flannel,/opt/cni/bin/flannel
./deploy/bin/host-device,/opt/cni/bin/host-device
./deploy/bin/portmap,/opt/cni/bin/portmap
./deploy/bin/sample,/opt/cni/bin/sample
./deploy/bin/vlan,/opt/cni/bin/vlan
./deploy/kubelet/10-weave.conf,/etc/cni/net.d/10-weave.conf
./deploy/ssl/etcd/ca.pem,/etc/etcd/ssl/ca.pem
./deploy/ssl/etcd/etcd.pem,/etc/etcd/ssl/etcd.pem
./deploy/ssl/etcd/etcd-key.pem,/etc/etcd/ssl/etcd-key.pem
./deploy/etcd/etcd3.service,/etc/systemd/system/etcd3.service
./deploy/etcd/etcd_ssl.sh,/opt/etcd_ssl.sh
./deploy/etcd/etcd_service.sh,/opt/etcd_service.sh
./scripts/resolve_ip.sh,/opt/resolve_ip.sh
./deploy/RestfulAPI/config.yaml,/etc/RestfulAPI/config.yaml
./deploy/master/restapi-kubeconfig.yaml,/etc/kubernetes/restapi-kubeconfig.yaml
./deploy/WebUI/userconfig.json,/etc/WebUI/userconfig.json
./deploy/WebUI/configAuth.json,/etc/WebUI/configAuth.json
./deploy/WebUI/dashboardConfig.json,/etc/WebUI/dashboardConfig.json
./deploy/WebUI/Master-Templates.json,/etc/WebUI/Master-Templates.json
./deploy/RestfulAPI/config.yaml,/etc/RestfulAPI/config.yaml
./deploy/storage/auto_share/auto_share.timer,/etc/systemd/system/auto_share.timer
./deploy/storage/auto_share/auto_share.target,/etc/systemd/system/auto_share.target
./deploy/storage/auto_share/auto_share.service,/etc/systemd/system/auto_share.service
./deploy/storage/auto_share/logging.yaml,{{cnf["folder_auto_share"]}}/logging.yaml
./deploy/storage/auto_share/auto_share_nfs_v2.py,{{cnf["folder_auto_share"]}}/auto_share.py
./deploy/storage/auto_share/mounting.yaml,{{cnf["folder_auto_share"]}}/mounting.yaml