# src/ClusterBootstrap/template/kubelet/deploy.list
./deploy/master/40-ExecStartPre-symlink.conf,/etc/systemd/system/flanneld.service.d/40-ExecStartPre-symlink.conf
./deploy/master/40-flannel.conf,/etc/systemd/system/docker.service.d/40-flannel.conf
./deploy/kubelet/options.env,/etc/flannel/options.env
./deploy/kubelet/kubelet.service,/etc/systemd/system/kubelet.service
./deploy/kubelet/worker-kubeconfig.yaml,/etc/kubernetes/worker-kubeconfig.yaml
./deploy/kubelet/dns-kubeconfig.yaml,/etc/kubernetes/dns-kubeconfig.yaml
./deploy/kubelet/kubelet.sh,/opt/kubelet.sh
./deploy/bin/kubelet,/opt/bin/kubelet
./deploy/ssl/ca/ca.pem,/etc/kubernetes/ssl/ca.pem
./deploy/ssl/kubelet/apiserver.pem,/etc/kubernetes/ssl/worker.pem
./deploy/ssl/kubelet/apiserver-key.pem,/etc/kubernetes/ssl/worker-key.pem
./deploy/kubelet/report.sh,/opt/report.sh
./deploy/kubelet/reportcluster.service,/etc/systemd/system/reportcluster.service
./deploy/kubelet/nodelist.yaml,/etc/kubernetes/nodes/nodelist.yaml
./deploy/kubelet/daemon.json,/etc/docker/daemon.json