# src/ClusterBootstrap/template/storage/auto_share/auto_share.target
[Unit]
Description=Auto Share Target
StopWhenUnneeded=yes